# fpmbot.timer
[Unit]
Description=Timer to start a fpmbot build

[Timer]
OnUnitActiveSec=6h
AccuracySec=1m
Persistent=true

[Install]
WantedBy=multi-user.target
