# fpmbot.service
[Unit]
Description=Execute fpmbot

[Service]
Type=oneshot
ExecStart=/bin/true

[Install]
WantedBy=multi-user.target
